# forwarder.service
[Unit]
Description=Telegram Channel Forwarder Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/telegram-forwarder
ExecStart=/usr/bin/python3 /opt/telegram-forwarder/main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
